# systemd/mfrontend.service
[Service]
WorkingDirectory=/opt/mirador/ui
ExecStart=/usr/bin/yarn start

[Install]
WantedBy=multi-user.target
