# systemd/mirador.service
[Service]
ExecStart=/opt/mirador/.venv/bin/python /opt/mirador/manage.py runserver

[Install]
WantedBy=multi-user.target
